# addons/pos_pot/rootfs/etc/services.d/pos/run
#!/usr/bin/with-contenv bash
# Start POS Pot addon via s6-overlay
exec /app/run.sh
